# supervisord/sms.ini
[program:sms]
command =  go run main.go
directory = /root/sms/
redirect_stderr = true
stdout_logfile = /data/logs/console/sms.log
autostart = true
autorestart = true
startsecs = 2
stopwaitsecs = 600
stdout_logfile_maxbytes=2MB 
stdout_logfile_backups=5
stderr_logfile_maxbytes=2MB 
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
environment = environment='test'
